# deploy/agent-ops-dashboard.init
#!/bin/sh
### BEGIN INIT INFO
# Provides:          agent-ops-dashboard
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Poetic agent-ops monitoring dashboard (loopback web server)
### END INIT INFO
#
# Install: see the "Run as a service" section of README.md. In short —
#   sudo install -m 755 deploy/agent-ops-dashboard.init /etc/init.d/agent-ops-dashboard
#   sudo service agent-ops-dashboard start
# and add it to /etc/wsl.conf's [boot] command, alongside cron.
#
# Serves scripts/serve-dashboard.sh's loopback-only web server (127.0.0.1
# only, never the network) so the dashboard loads over http:// when a
# browser refuses data.js over a file:// URL. Runs as the unprivileged
# `wallen` user (never root) via start-stop-daemon --chuid. HOME and PATH
# are set explicitly because this script itself runs as root (via WSL's
# [boot] command) — left alone, $HOME would stay /root through the chuid
# and break config.json's ~-expansion and the dashboard's state-dir lookup,
# and a bare boot PATH may not include jq/python3.
#
# Every setting below is a default that /etc/default/agent-ops-dashboard may
# override, so a host with a different user, home or checkout path needs no
# edit to this file. This is the legacy single-laptop path; a containerised
# node runs serve-dashboard.sh directly and lets the container runtime
# capture its output (see deploy/docker/).

# shellcheck disable=SC1091  # optional host-local override file
[ -r /etc/default/agent-ops-dashboard ] && . /etc/default/agent-ops-dashboard

NAME=agent-ops-dashboard
RUNAS=${RUNAS:-wallen}
RUNHOME=${RUNHOME:-/home/$RUNAS}
APPDIR=${APPDIR:-$RUNHOME/Code/Poetic-Poems/agent-ops}
PORT=${PORT:-8787}
PIDFILE=${PIDFILE:-/var/run/$NAME.pid}
PATH_ENV=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:$RUNHOME/.local/bin
# The server's own log belongs with the rest of the pipeline's state, not
# beside the checkout: state_dir is the one directory a node backs up and a
# container mounts. Read it from config.json (expanding a leading ~ against
# the run user's home, as the shell scripts do) and fall back to the
# documented default if jq or the config is unreadable at boot. PATH is
# forced for the lookup because root's boot PATH may not include jq.
if [ -z "$LOGFILE" ]; then
  state_dir=$(PATH=$PATH_ENV jq -r '.state_dir' "$APPDIR/config.json" 2>/dev/null || true)
  case "$state_dir" in
    "~"*) state_dir=$RUNHOME${state_dir#\~} ;;
    ""|null) state_dir=$RUNHOME/.local/state/poetic-agents ;;
  esac
  LOGFILE=$state_dir/dashboard-server.log
fi
CMD="HOME=$RUNHOME; export HOME; PATH=$PATH_ENV; export PATH; mkdir -p \"\$(dirname $LOGFILE)\"; exec $APPDIR/scripts/serve-dashboard.sh $PORT >>$LOGFILE 2>&1"

case "$1" in
  start)
    # start-stop-daemon can't reliably match this against a running process
    # by --exec: the tracked process execs from /bin/bash into python3, so
    # its comm name no longer matches what --start would look for. Guard
    # against a duplicate server (a second python http.server on the same
    # port dies with "Address already in use") with an explicit pidfile
    # check instead.
    if [ -f "$PIDFILE" ]; then
      OLDPID=$(cat "$PIDFILE" 2>/dev/null)
      if [ -n "$OLDPID" ] && kill -0 "$OLDPID" 2>/dev/null; then
        echo "$NAME is already running (pid $OLDPID)"
        exit 0
      fi
      echo "Removing stale pidfile for $NAME"
      rm -f "$PIDFILE"
    fi
    echo "Starting $NAME on http://127.0.0.1:$PORT"
    start-stop-daemon --start --background --make-pidfile --pidfile "$PIDFILE" \
      --chuid "$RUNAS:$RUNAS" --chdir "$APPDIR" \
      --exec /bin/bash -- -c "$CMD"
    ;;
  stop)
    echo "Stopping $NAME"
    start-stop-daemon --stop --pidfile "$PIDFILE" --retry TERM/10/KILL/5
    rm -f "$PIDFILE"
    ;;
  restart)
    "$0" stop
    sleep 1
    "$0" start
    ;;
  status)
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
      echo "$NAME is running (pid $(cat "$PIDFILE"))"
    else
      echo "$NAME is not running"
      exit 1
    fi
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac
